# data/hyperstart.service.in
[Unit]
Description=hyperstart guest daemon

[Service]
ExecStart=@prefix@/bin/hyperstart
